# pkg/cloudflared.rc
#!/bin/ksh

daemon="${PREFIX}/bin/cloudflared"
daemon_user="nobody"
daemon_flags="tunnel run"

. /etc/rc.d/rc.subr

rc_bg=YES
rc_reload=NO

rc_cmd $1
